# etc/supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/
user=root
autostart=true
autorestart=true
